# docker/supvervisord/caddy.ini
[program:caddy]
command=/usr/sbin/caddy run --config /etc/caddy/Caddyfile
