# docker/supvervisord/php-fpm.ini
[program:php-fpm]
command=/usr/sbin/php-fpm82 -F
autostart=true
autorestart=unexpected
stdout_logfile=/var/log/stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/stderr.log
stderr_logfile_maxbytes=0
exitcodes=0
